# intellirouter-router.service
[Unit]
Description=IntelliRouter Router Service
After=network.target
Wants=redis-server.service

[Service]
Type=simple
User=intellirouter
Group=intellirouter
WorkingDirectory=/opt/intellirouter
EnvironmentFile=/opt/intellirouter/env/router.env
ExecStart=/opt/intellirouter/bin/intellirouter run --role Router --config /opt/intellirouter/config/router.toml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
StandardOutput=append:/opt/intellirouter/logs/router/stdout.log
StandardError=append:/opt/intellirouter/logs/router/stderr.log

# Security settings
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target